# frontend/scripts/electro-next.service
[Unit]
Description=Next.js Garant Montazh
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/electro-garant/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/env node ./node_modules/next/dist/bin/next start -H 0.0.0.0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
